# php/daemon/build/supervisord.conf
[supervisord]
nodaemon=true

[program:process]
command=php /app/console Files:Process teste